# workspace/crontab/laradock
#* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock cd /var/www/emop-api.com && php artisan schedule:run >> /dev/null 2>&1
